# loaduvc.service
[Unit]
Description=loads custom uvc driver

[Service]
Type=oneshot
ExecStart=/usr/bin/loaduvc.sh

[Install]
WantedBy=multi-user.target