# loaduvc.service
[Unit]
Description=loads custom uvc driver

[Service]
Type=oneshot
ExecStart=/root/modules/loaduvc.sh

[Install]
WantedBy=multi-user.target